# systemd/nano-snapshot.service
[Unit]
Description=Nano P2P Snapshot Pipeline (download, publish, kick seeder)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-%h/.env
Environment=OUTPUT_DIR=%h/nano-snapshots
Environment=WEB_SEED_URL=https://s3.us-east-2.amazonaws.com/repo.nano.org/snapshots/latest
# Adjust ExecStart if the repo is cloned to a different location
ExecStart=%h/nano-snapshot-swarm/scripts/daily-snapshot.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nano-snapshot
# 57GB download + torrent hashing + DHT publish can take many hours
TimeoutStartSec=43200
TimeoutStopSec=300
