# systemd/nano-status-push.timer
[Unit]
Description=Push snapshot status to status API hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
